# systemd/check-go-collector.service
[Unit]
Description=Check go-collector health and self-heal
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu

EnvironmentFile=/etc/ubuntu-wallet/telegram.env
ExecStart=/usr/bin/env bash -lc '/home/ubuntu/ubuntu-wallet/scripts/ops/check-go-collector.sh >> /home/ubuntu/ubuntu-wallet/check-go-collector.log 2>&1'
