# systemd/daily-report.service
[Unit]
Description=ubuntu-wallet daily evaluation report generator (per-symbol from configs/symbols.yaml)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/ubuntu-wallet
ExecStartPre=/bin/mkdir -p /home/ubuntu/ubuntu-wallet/data/logs

# Strict YAML mode:
# - Uses --all-symbols, and requires threshold/tp/sl/horizon/mt_filter_mode be present per symbol in configs/symbols.yaml
# - Writes reports to data/<SYMBOL>/reports/ by default (do NOT pass --report-dir)
ExecStart=/bin/bash -c 'source /home/ubuntu/ubuntu-wallet/ml-service/.venv/bin/activate && \
  python /home/ubuntu/ubuntu-wallet/scripts/generate_daily_report.py \
    --all-symbols \
    --active-model event_v3 \
    >> /home/ubuntu/ubuntu-wallet/data/logs/daily_report.log 2>&1'
