# systemd/evaluate-predictions.timer
[Unit]
Description=Run prediction evaluator every 6 hours
Requires=evaluate-predictions.service

[Timer]
# Run 6 minutes past the hour, every 6 hours
OnCalendar=*-*-* 00,06,12,18:06:00
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
